# files/etc/supervisor.d/controller.ini
[program:controller]
command=/controller.sh
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=true
environment=MAX_CONN=%(ENV_MAX_CONN)s
